# seven-control/location/seven-control-location-agent.service
[Unit]
Description=Seven Control Location Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/seven-control/location-agent.env
ExecStart=/usr/local/bin/seven-control-location-agent
User=root
